# files/profile.d/010-clamav.sh
#!/usr/bin/env bash

PATH="${PATH}:${HOME}/clamav/bin:${HOME}/clamav/sbin"
export PATH

LD_LIBRARY_PATH="${LD_LIBRARY_PATH}:${HOME}/clamav/lib"
export LD_LIBRARY_PATH
