# files/profile.d/010-clamav.sh
#!/usr/bin/env bash

PATH="${PATH}:${HOME}/bin:${HOME}/sbin"
export PATH

LD_LIBRARY_PATH="${LD_LIBRARY_PATH}:${HOME}/lib"
export LD_LIBRARY_PATH
